# common/system/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin badele --noclear %I 38400 linux
Type=simple
